# config/fb-unipdbot.conf
[program:fb-unipdbot]
command = gunicorn app:app -b localhost:8000
directory = /home/michele/bot_app
user = michele